# nodejs/phantom_proxy.conf
start on started

script
  for env in $(cat /proc/1/environ | tr "\000" "\n"); do
    export $env
  done
  unset http_proxy
  node /phantom_proxy.js 2>&1 | logger -i -t phantom_proxy
end script
